# core/build/conf/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/run/supervisord.pid
# On évite que supervisord ne garde trop de logs internes
logfile_maxbytes=50MB
logfile_backups=10

[program:php-fpm]
# Le flag -R permet de lancer en root (supervisord s'en occupe) 
# mais FPM switchera sur 'bankuser' comme défini dans www.conf
command=/usr/sbin/php-fpm -F -O
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=5
# TRÈS IMPORTANT : Force PHP à envoyer les erreurs au worker
# qui les renverra ensuite à supervisord -> Docker logs
environment=PHP_INI_SCAN_DIR=:/etc/php.d

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=10
# Prévient Apache de ne pas créer de fichiers PID inutiles qui bloquent le restart
exec=true